# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -p -m 30
